# backup2nas.service
[Unit]
Description=Backup to NAS
Requires=mnt-backup.mount

[Service]
Type=oneshot
ExecStart=/usr/local/bin/backup2nas

